# allcanlearn.service
[Unit]
Description=AllCanLearn App
After=network.target

[Service]
Type=simple
User=vikki
WorkingDirectory=/home/vikki/AllCanLearn-v6
Environment=PATH=/home/vikki/AllCanLearn-v6/venv/bin
EnvironmentFile=/home/vikki/AllCanLearn-v6/.env
ExecStart=/home/vikki/AllCanLearn-v6/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=allcanlearn

[Install]
WantedBy=multi-user.target
